# ruby-test/tailscale-up-if-authkey.service
[Unit]
Description=Authenticate Tailscale when TS_AUTHKEY is set
After=tailscaled.service
Requires=tailscaled.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/tailscaled-auth
PassEnvironment=TS_AUTHKEY TS_EXTRA_ARGS
ExecStart=/usr/local/bin/tailscale-up-if-authkey
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
